# obarun-boot-0.1/services/dhcpcd/run
#!/bin/sh

[ -r conf ] && . ./conf

exec /usr/bin/dhcpcd -B enp2s0 $OPTS 1>&2


